# templates/11-mjpg_streamer.rules
# file /etc/udev/rules.d/11-mjpgstreamer.rules
KERNEL=="video?", ACTION=="add", RUN+="/usr/sbin/service mjpg_streamer restart"
